# yolo-dev.service
[Unit]
Description=YOLO Dev Inference API Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/PycharmProjects/YoloService
ExecStart=/home/deema/PycharmProjects/YoloService/venv/bin/python3 /home/deema/PycharmProjects/YoloService/app.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

