# yolo-dev.service
[Unit]
Description=YOLO Dev Inference API Service
After=network.target

[Service]
WorkingDirectory=/home/ubuntu/yolo-dev
ExecStart=/home/ubuntu/yolo-dev/venv/bin/python3 /home/ubuntu/yolo-dev/app.py
Restart=always
User=ubuntu
Environment=ENV=development

[Install]
WantedBy=multi-user.target
